# rclone-home/homecp/.config/systemd/user/rclone-home.timer
[Unit]
Description=Run rclone-home daily

[Timer]
OnBootSec=15min
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target